# docker_files/docker.conf
[supervisord]
nodaemon=true

[inet_http_server]
port: :9999
username: docker
password: d0ck3r

[program:redis]
command=/usr/bin/redis-server
redirect_stderr=true
autorestart=true

